# files/cassandra.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: cassandra
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name=cassandra
rcvar=cassandra_enable

load_rc_config $name

: ${cassandra_enable:=NO}
: ${cassandra_pidfile=/var/run/cassandra.pid}

command=%%PREFIX%%/bin/cassandra
command_args="-p ${cassandra_pidfile}"

start_precmd="export JAVA_HOME=%%JAVA_HOME%%"

procname=%%PREFIX%%/openjdk8/bin/java
pidfile=${cassandra_pidfile}

run_rc_command "$1"
